# roles/geoipupdate/files/etc/systemd/system/geoipupdate.service
[Unit]
Description=MaxMind GeoIP/GeoIP2 database updates

[Service]
Type=oneshot
ExecStart=/usr/bin/geoipupdate

[Install]
WantedBy=default.target
