# roles/iptables/files/etc/network/if-pre-up.d/iptables
#!/bin/sh

# Ansible managed.

/sbin/iptables-legacy-restore < /etc/iptables.up.rules
