# roles/dehydrated/files/etc/cron.daily/dehydrated
#!/bin/sh
/usr/bin/dehydrated -c > /dev/null
